# scripts/systemd/fikrat-mail.service
[Unit]
Description=Fikrat: Experts d'assurance (MAIL)
After=network.target

[Service]
WorkingDirectory=/opt/fikrat
; Environment=DEBUG=app:*
ExecStart=/usr/bin/yarnpkg run mail
; Restart=always

[Install]
WantedBy=multi-user.target
